# sshd.conf
[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
exitcodes=0
startretries=999
priority=1
